# systemd/reprepro-export@.timer
[Unit]

[Timer]
OnActiveSec=0
OnUnitInactiveSec=48h

[Install]
WantedBy=timers.target
