# systemd/reprepro-update@.timer
[Unit]

[Timer]
OnActiveSec=0
OnUnitInactiveSec=8h

[Install]
WantedBy=timers.target
